# ads/opctl/docker/etc/supervisord.d/jupyterlab.ini
[program:jupyterlab]
directory=/home/datascience
command=/etc/supervisord.d/jupyterlab/run
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=3
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=%(ENV_JUPYTER_SERVER_LOG_DIR)s/%(ENV_JUPYTER_SERVER_LOG_FILE)s
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
stopsignal=TERM
stopwaitsecs=60
redirect_stderr=true