# ads/opctl/docker/etc/supervisord.d/nbruntime.ini
[program:nbruntime]
directory=/home/datascience
command=/bin/bash -c "exec /etc/supervisord.d/nbruntime/run"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=3
process_name=%(ENV_NBRUNTIME_SERVICE_NAME)s
stdout_logfile=%(ENV_JUPYTER_SERVER_LOG_DIR)s/%(ENV_NBRUNTIME_SERVICE_LOG_FILE)s
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
stopsignal=TERM
stopwaitsecs=%(ENV_MAX_SCRIPT_EXECUTION_SECONDS)s
redirect_stderr=true